# services/dns.timer
[Unit]
Description=Update godaddy dns

[Timer]
OnBootSec=1min
OnCalendar=minutely

[Install]
WantedBy=timers.target
